# images/proton-bridge/root/etc/services.d/proton-bridge-imap/run
#!/command/with-contenv bash

exec /usr/bin/socat TCP-LISTEN:143,fork TCP:localhost:1143
