# images/restic/root/etc/crontabs/root
{{ RESTIC_SCHEDULE }} /app/restic-backup.sh >> /config/log/restic/restic.log 2>&1
30 5 1 * * /app/restic-check.sh >> /config/log/restic/restic.log 2>&1
