# images/restic/root/etc/crontabs/root
{{ RESTIC_SCHEDULE }} s6-setuidgid abc /app/restic-backup.sh >> /config/log/restic/restic.log 2>&1
30 5 1 * * s6-setuidgid abc /app/restic-check.sh >> /config/log/restic/restic.log 2>&1
